# .config/systemd/user/clash.service
[Unit]
Description=Clash autostarter
After=network.target

[Service]
Type=simple
ExecStart=/home/peviroy/Projects/MyBash/clash.sh
StandardOutput=null
StandardError=journal
Restart=on-abort

[Install]
WantedBy=multi-user.target
